# srv/immich/immich.service
[Unit]
Description=Immich Stack
Wants=network-online.target
After=network-online.target local-fs.target

[Service]
Type=oneshot
RemainAfterExit=yes

WorkingDirectory=%h/homelab/srv/immich
ExecStart=/usr/bin/podman compose up -d
ExecStop=/usr/bin/podman compose down
ExecReload=/usr/bin/podman compose restart

[Install]
WantedBy=default.target
